# smarthome-display.service
[Unit]
Description=SelenaCore — TTY1 display (wizard QR / system stats)
After=network.target docker.service
Requires=docker.service
Conflicts=getty@tty1.service

[Service]
Type=simple
User=root

# Wait until the core container is healthy (search by compose service label)
ExecStartPre=/bin/bash -c \
  'for i in $(seq 1 30); do \
     docker compose -f /home/selena/SelenaCore/docker-compose.yml ps core 2>/dev/null | grep -qE "running|healthy" && exit 0; \
     echo "[selena-display] waiting for core service... ($i)"; \
     sleep 3; \
   done; echo "[selena-display] timeout"; exit 1'

# Run tty_status.py via docker compose exec; pass the host real IP for the QR code URL
ExecStart=/bin/bash -c \
  'HOST_IP=$(hostname -I | awk "{print $1}"); \
   docker compose -f /home/selena/SelenaCore/docker-compose.yml \
     exec -T -e HOST_IP="$HOST_IP" core \
     python -m system_modules.ui_core.tty_status > /dev/tty1 2>&1'

Restart=always
RestartSec=5
TimeoutStartSec=120

StandardOutput=journal
StandardError=journal
SyslogIdentifier=selena-display

[Install]
WantedBy=multi-user.target
