# smarthome-modules.service
[Unit]
Description=SelenaCore — Modules Container Runtime
After=docker.service smarthome-core.service
Requires=docker.service

[Service]
Type=simple
User=selena
Group=selena
ExecStart=/usr/bin/docker run --rm \
    --name selena-modules \
    --network selena_selena_external \
    -p 8100-8200:8100-8200 \
    ${MODULE_CONTAINER_IMAGE}
ExecStop=/usr/bin/docker stop selena-modules
Restart=always
RestartSec=10

EnvironmentFile=/opt/selena-core/.env

[Install]
WantedBy=multi-user.target
