# whisper-server.service
[Unit]
Description=Whisper.cpp STT Server (CUDA)
After=network.target
Before=smarthome-core.service

[Service]
Type=simple
User=root
EnvironmentFile=-/var/lib/selena/whisper-server.env
Environment=LD_LIBRARY_PATH=/home/dotradepro/whisper.cpp/build/ggml/src/ggml-cuda:/usr/local/cuda-12.6/lib64
Environment=WHISPER_MODEL=ggml-small
ExecStart=/home/dotradepro/whisper.cpp/build/bin/whisper-server \
    --model /home/dotradepro/whisper.cpp/models/${WHISPER_MODEL}.bin \
    --host 0.0.0.0 \
    --port 9000 \
    --language auto
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
