# setup/pulseaudio-system.service
[Unit]
Description=PulseAudio sound server (system mode for headless kiosk)
Before=getty@tty1.service

[Service]
Type=simple
User=dotradepro
Group=dotradepro
SupplementaryGroups=audio pulse
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStart=/usr/bin/pulseaudio --disallow-exit --disallow-module-loading=0 --exit-idle-time=-1 --log-target=journal
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
